# systemd/ntm-notify-monitor.service
[Unit]
Description=NTM Agent Notification Monitor
After=network.target

[Service]
Type=simple
ExecStart=%h/.local/bin/ntm-notify-monitor.sh
Restart=on-failure
RestartSec=5
Environment=HOME=%h
Environment=PATH=%h/.local/bin:/usr/local/bin:/usr/bin:/bin

[Install]
WantedBy=default.target
